# deploy/systemd/smm-worker.service
[Unit]
Description=SMM Panel Background Worker
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=smm
WorkingDirectory=/opt/smm
EnvironmentFile=/opt/smm/.env
ExecStart=/opt/smm/bin/worker
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
