# deploy/systemd/smm-server.service
[Unit]
Description=SMM Panel Webhook Server
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=smm
WorkingDirectory=/opt/smm
EnvironmentFile=/opt/smm/.env
ExecStart=/opt/smm/bin/server
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
